# autosyncoid.timer

[Unit]
Description=Runs autosyncoid script

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
